# docker/rootfs/etc/cron.d/www-data
* * * * * curl http://localhost/instance-deleteExpiredDemoInstance.html
