# docker/rootfs/etc/cron.d/www-data
PATH=/usr/local/bin:/usr/local/sbin:/usr/bin:/usr/sbin:/bin:/sbin

* * * * * curl http://localhost/instance-deleteExpiredDemoInstance.html
* * * * * /bin/sh /apps/qucheng/bin/backup.sh
